# powertopUSB.service
[Unit]
Description=Powertop USB tunings

[Service]
ExecStart=

[Install]
WantedBy=multi-user.target
