# powertopUSB.service
[Unit]
Description=Powertop USB tunings

[Service]
ExecStart= /usr/bin/bash /usr/bin/powertopUSB

[Install]
WantedBy=multi-user.target
