# powertopUSB.service
[Unit]
Description=Powertop USB tunings

[Service]
ExecStart= ~/.Scripts

[Install]
WantedBy=multi-user.target
